# assets/marxist-quote-fetch.service
[Unit]
Description=Fetch daily quote from WikiQuote
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/marxist-quote-background --fetch
ExecStartPost=/usr/bin/systemctl --user restart desktop-quote.service
